# systemd/cowrie-split.path
[Unit]
Description=Watch Cowrie JSONL log for changes

[Path]
PathModified=/home/cowrie/honeypot/logs/raw-logs/cowrie.json
TriggerLimitIntervalSec=1s

[Install]
WantedBy=multi-user.target
